# scripts/m5-voice/stack/systemd/mango-companion.service
[Unit]
Description=mango voice companion HTTPS (phone PTT)
After=network-online.target mango-orchestrator.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mango
Environment=HOME=%h
Environment=MANGO_SKIP_OVERLAY=1
Environment=MANGO_COMPANION_HOST=0.0.0.0
Environment=MANGO_COMPANION_PORT=3001
ExecStartPre=/bin/mkdir -p %h/.cache/mango
ExecStartPre=/bin/bash -c 'test -d %h/mango/src/companion/dist || npm --prefix %h/mango/src/companion run build'
ExecStart=/bin/bash -lc 'source %h/.config/mango/voice.env 2>/dev/null || true; exec bash %h/mango/scripts/m5-voice/stack/serve-companion-https.sh'
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.cache/mango/companion.log
StandardError=append:%h/.cache/mango/companion.log

[Install]
WantedBy=default.target
