# config/systemd/mango-aiolists.service
[Unit]
Description=mango self-hosted AIOLists
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mango/deploy/aiolists
Environment=HOME=%h
ExecStart=/usr/bin/env docker compose up
ExecStop=/usr/bin/env docker compose down
Restart=on-failure
RestartSec=10
TimeoutStartSec=180
TimeoutStopSec=30

[Install]
WantedBy=default.target
