# scripts/m1-foundation/ui/systemd/mango-watchdog.service
[Unit]
Description=mango TV health watchdog
Wants=mango-ui-server.service mango-catalog.service mango-launcher-chromium.service mango-tv-pad.service
After=mango-ui-server.service mango-catalog.service mango-launcher-chromium.service mango-tv-pad.service

[Service]
Type=oneshot
Environment=HOME=%h
Environment=MANGO_REPO_DIR=%h/mango
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=%h/mango/scripts/mango-health-repair.sh --quiet
