# scripts/m1-foundation/ui/systemd/mango-catalog.service
[Unit]
Description=mango catalog-service
After=network-online.target mango-aiometadata.service
Wants=network-online.target mango-aiometadata.service

[Service]
Type=simple
WorkingDirectory=%h/mango
Environment=HOME=%h
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=MANGO_CATALOG=1
ExecStartPre=/bin/mkdir -p %h/.cache/mango
ExecStart=%h/mango/scripts/m2-catalog/service/run-catalog-service.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=180

[Install]
WantedBy=default.target
