# scripts/phase2/systemd/mango-orchestrator.service
[Unit]
Description=mango voice orchestrator (WSS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mango
Environment=HOME=%h
Environment=MANGO_SKIP_OVERLAY=1
Environment=MANGO_ORCH_TLS=1
ExecStartPre=/bin/mkdir -p %h/.cache/mango
ExecStartPre=%h/mango/scripts/phase2/ensure-orchestrator-venv.sh
ExecStart=/bin/bash -lc 'source %h/.config/mango/voice.env 2>/dev/null || true; cd %h/mango && exec bash scripts/phase2/start-orchestrator.sh'
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.cache/mango/orchestrator.log
StandardError=append:%h/.cache/mango/orchestrator.log

[Install]
WantedBy=default.target
